# service/virtualbox-restore@.service
[Unit]
Description=VirtualBox restore service
After=sleep.target

[Service]
User=%I
Type=forking
ExecStart=/home/%I/bin/virtualbox-ctl restore

[Install]
WantedBy=sleep.target